# scripts/podman/openclaw.container.in
# gensparx gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into gensparx's ~/.config/containers/systemd/
# {{GENSPARX_HOME}} is replaced at install time.

[Unit]
Description=gensparx gateway (rootless Podman)

[Container]
Image=gensparx:local
ContainerName=gensparx
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{GENSPARX_HOME}}/.gensparx:/home/node/.gensparx
EnvironmentFile={{GENSPARX_HOME}}/.gensparx/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
