# scripts/systemd/openclaw-auth-monitor.service
[Unit]
Description=gensparx Auth Expiry Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/home/admin/gensparx/scripts/auth-monitor.sh
# Configure notification channels via environment
Environment=WARN_HOURS=2
# Environment=NOTIFY_PHONE=+1234567890
# Environment=NOTIFY_NTFY=gensparx-alerts

[Install]
WantedBy=default.target
